# src/modules/beamos/filesystem/mount_manager/etc/systemd/system/usb_mount_manager_remove.service
[Service]
Type=simple
ExecStart=/root/mount_manager/mount_manager remove
