# src/modules/beamos/filesystem/mount_manager/root/etc/systemd/system/mount_manager_remove.service
[Service]
Type=simple
ExecStart=/usr/bin/mount_manager remove

[Install]
WantedBy=mount_manager_add.service
