# src/modules/beamos/filesystem/mount_manager/root/etc/systemd/system/usb_mount_manager_remove.service
[Service]
Type=simple
ExecStart=sudo -u pi /usr/bin/mount_manager remove
